# Idea/backup-database/backup.cron
0 6 * * * sh /home/data/mysqlDump/dump.sh
